# etc/supervisord.conf

[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:app]
command=bin/supervisord -c /app/root/etc/docker-supervisord.conf -n
directory=/app/root
autostart=false
